# billing/etc/cron/crontab
0 0 * * * python manage.py process_billing_cycle
